# xg2g.service
[Unit]
Description=xg2g - Enigma2 to IPTV Gateway
Documentation=https://github.com/ManuGH/xg2g
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/xg2g
EnvironmentFile=/root/xg2g/.env
Environment=XG2G_DATA=/root/xg2g/data

# Start xg2g daemon
ExecStart=/root/xg2g/bin/xg2g

# Restart configuration
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal

# Security & Limits
# 1. NoNewPrivileges=true: Ensures that the service and its children can never gain new privileges (e.g. via setuid).
# 2. PrivateTmp=true: Mounts a private /tmp for the service, isolating it from the host's /tmp.
# 3. LimitNOFILE=65535: Increases the file descriptor limit to handle high-concurrency streaming.
# 4. UMask=0077: Ensures files created by the service are only readable by the owner (root), protecting secrets/config.
NoNewPrivileges=true
PrivateTmp=true
LimitNOFILE=65535
UMask=0077

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
