# backend/cmd/daemon/deploy/xg2g.service
[Unit]
Description=xg2g
[Service]
EnvironmentFile=/etc/xg2g/xg2g.env
